# launchagent/com.daily-report.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.__YOUR_USERNAME__.daily-report</string>
  <key>ProgramArguments</key>
  <array>
    <string>__HOME__/bin/run-daily-report.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>21</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/daily-report.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/daily-report.err.log</string>
</dict>
</plist>
